# data/whiteknight.service.in
[Unit]
Description=Whitelist driven mandatory access control application
Requires=@requires@

[Service]
ExecStart=@bindir@/whiteknight

[Install]
WantedBy=sysinit.target
